# dial.service
[Unit]
Description=dial

[Service]
ExecStart=/home/pi/vibtest/dial.sh
RemainAfterExit=true
[Install]

WantedBy=multi-user.target


